# scripts/broute.service
[Unit]
Description = B route data fetch service

[Service]
Type = oneshot
ExecStart = /opt/scripts/b-route-reader/scripts/broute.sh
